# packaging/systemd/bind9-api.service
[Unit]
Description=Bind9 API Service
After=network.target

[Service]
ExecStart=/usr/local/bin/bind9-api
Restart=always
RestartSec=5
User=nobody
Group=nogroup
Environment=GIN_MODE=release

[Install]
WantedBy=multi-user.target
